# buildroot-external/board/rpi4/overlay/usr/lib/systemd/system/aeroscan-gui.service
[Unit]
Description=AeroScan GUI
After=sysinit.target dump1090.service
# dump1090 runs only alongside the GUI (see PartOf= in dump1090.service);
# its own 5 s ExecStartPre covers USB enumeration, and the GUI's ADS-B
# worker reconnects to port 30003 every 5 s, so ordering is best-effort.
Wants=dump1090.service

[Service]
Type=simple
User=root
Environment=QT_FONT_DPI=84
Environment=QT_QPA_PLATFORM=eglfs
# Static defaults for the Waveshare HDMI panel (480x800 portrait: render
# rotated 90° — EGLFS_ROTATION only rotates rendering, never input).
# Touch (XPT2046): panel axes are swapped AND both inverted vs the landscape
# view; inversions are applied before rotation on normalized evdevtouch
# coords, hence rotate=270:invertx:inverty (verified empirically 2026-07-03).
Environment=QT_QPA_EGLFS_KMS_CONFIG=/etc/qt-kms.json
Environment=QT_QPA_EGLFS_ROTATION=90
Environment=QT_QPA_EVDEV_TOUCHSCREEN_PARAMETERS=/dev/input/touchscreen0:rotate=270:invertx:inverty
# aeroscan-display-init picks the attached display (DSI preferred) and writes
# display.env; values from EnvironmentFile override the static defaults above.
ExecStartPre=/usr/sbin/aeroscan-display-init
EnvironmentFile=-/run/aeroscan/display.env
ExecStart=/usr/bin/aeroscan-gui
# Appliance mode: the GUI must outlive any crash, indefinitely. Power-menu
# exits are deliberate and must NOT respawn: 41=poweroff, 42=reboot,
# 43=exit-to-terminal (codes from winglet-ui/theme.h; actions performed by
# aeroscan-power-action below). Exit 0 = Restart UI — respawns by design.
Restart=always
RestartSec=3
StartLimitIntervalSec=0
SuccessExitStatus=41 42 43
RestartPreventExitStatus=41 42 43
ExecStopPost=/usr/sbin/aeroscan-power-action
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
