# buildroot-external/board/rpi4/overlay/usr/lib/systemd/system/dump1090.service
[Unit]
Description=dump1090 ADS-B decoder
After=network.target systemd-udev-settle.service
Wants=network.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/dump1090 \
    --net \
    --net-sbs-port 30003 \
    --net-ro-port 30002 \
    --quiet \
    --no-fix
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
